# apps/pi-prototype-cli/systemd/baddream-button.service.example
[Unit]
Description=Bad Dream Button Pi Runtime
After=network-online.target baddream-bridge-tunnel.service
Wants=network-online.target baddream-bridge-tunnel.service

[Service]
Type=simple
User=saleem
WorkingDirectory=/home/saleem/baddreambutton/apps/pi-prototype-cli
Environment=BAD_DREAM_CONFIG=/home/saleem/.config/baddream-button/config.json
ExecStart=/usr/bin/python3 /home/saleem/baddreambutton/apps/pi-prototype-cli/main.py run
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
